# driver-kafka/deploy/templates/zookeeper.service
[Unit]
Description=ZooKeeper
After=network.target

[Service]
ExecStart=/opt/kafka/bin/zookeeper-server-start.sh config/zookeeper.properties
WorkingDirectory=/opt/kafka
RestartSec=1s
Restart=on-failure
Type=simple

[Install]
WantedBy=multi-user.target
