# configuration/templates/pulsar-dashboard.service
[Unit]
Description=Pulsar Dashboard
After=docker.service
Requires=docker.service
After=prometheus.service
Requires=prometheus.service
 
[Service]
WorkingDirectory=/opt/grafana
ExecStart=/usr/libexec/docker/cli-plugins/docker-compose up
ExecStop=/usr/libexec/docker/cli-plugins/docker-compose down
TimeoutStartSec=0
Restart=on-failure
StartLimitBurst=3
 
[Install]
WantedBy=multi-user.target
