# driver-nats/deploy/templates/nats-server.service
[Unit]
Description=nats broker
After=network.target

[Service]
ExecStart=/opt/nats/gnatsd-v1.3.0-linux-amd64/gnatsd -c /opt/nats/cluster.conf
WorkingDirectory=/opt/nats/
RestartSec=1s
Restart=on-failure
Type=simple
LimitNOFILE=300000

[Install]
WantedBy=multi-user.target