# driver-pravega/deploy/templates/benchmark-worker.service
[Unit]
Description=OpenMessaging Benchmark Worker
After=network.target

[Service]
ExecStart=/opt/benchmark/bin/benchmark-worker
WorkingDirectory=/opt/benchmark
RestartSec=15s
Restart=on-failure
Type=simple
LimitNOFILE=300000

[Install]
WantedBy=multi-user.target
