# driver-pravega/deploy/templates/pravega-segmentstore.service
[Unit]
Description=Pravega Segment Store
After=network.target

[Service]
Environment="PRAVEGA_SEGMENTSTORE_OPTS=-Xms6G -Xmx6G -XX:MaxDirectMemorySize=24g"
ExecStart=/opt/pravega/bin/pravega-segmentstore
WorkingDirectory=/opt/pravega
RestartSec=15s
Restart=on-failure
Type=simple
LimitNOFILE=300000

[Install]
WantedBy=multi-user.target
