# driver-rabbitmq/deploy/templates/prometheus.service
[Unit]
Description=Prometheus server
After=network.target

[Service]
WorkingDirectory=/opt/prometheus
ExecStart=/opt/prometheus/prometheus \
    --web.enable-admin-api \
    --storage.tsdb.path=/opt/prometheus/data \
    --config.file=/opt/prometheus/prometheus.yml

[Install]
WantedBy=multi-user.target
