# driver-rocketmq/deploy/templates/rmq-namesrv.service
[Unit]
Description=RocketMQ-NameServer
After=network.target

[Service]
ExecStart=/opt/rocketmq/rocketmq-all-4.3.0-bin-release/bin/mqnamesrv -c conf/rmq-namesrv.properties
WorkingDirectory=/opt/rocketmq/rocketmq-all-4.3.0-bin-release/
RestartSec=1s
Restart=on-failure
Type=simple
LimitNOFILE=300000

[Install]
WantedBy=multi-user.target
