# driver-nats-streaming/deploy/templates/nats-streaming-server.service
[Unit]
Description=nats broker
After=network.target

[Service]
ExecStart=/opt/nats-streaming/nats-streaming-server-v0.25.2-linux-amd64/nats-streaming-server -c /opt/nats-streaming/cluster.conf
WorkingDirectory=/opt/nats-streaming/
RestartSec=1s
Restart=on-failure
Type=simple
LimitNOFILE=300000

[Install]
WantedBy=multi-user.target